# image/hack_back.conf
start on startup
task
exec /hack_back_ckpt.rcS
