# image/hack_back.conf
start on startup
task
exec /hack_back_impica_ckpt.rcS
